# scripts/launchd/com.brainlayer.enrichment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brainlayer.enrichment</string>

    <!-- Realtime Gemini enrichment LaunchAgent.
         Runs once at login/load and then hourly. -->

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>-c</string>
        <string>import sys; sys.path.insert(0, "__BRAINLAYER_DIR__/src"); from brainlayer.paths import get_db_path; from brainlayer.vector_store import VectorStore; from brainlayer.enrichment_controller import enrich_realtime; store = VectorStore(get_db_path()); result = enrich_realtime(store=store, limit=50, since_hours=24); print(result); store.close()</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__BRAINLAYER_DIR__</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/brainlayer-enrichment.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/brainlayer-enrichment.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:__HOME__/.local/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONPATH</key>
        <string>__BRAINLAYER_DIR__/src</string>
        <key>BRAINLAYER_STALL_TIMEOUT</key>
        <string>300</string>
        <key>GOOGLE_API_KEY</key>
        <string>__GOOGLE_API_KEY__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
