# scripts/launchd/com.brainlayer.enrichment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brainlayer.enrichment</string>

    <!-- Realtime Gemini enrichment LaunchAgent.
         Runs a supervised shell loop so the one-shot enrich command stays continuous. -->

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>cd "__BRAINLAYER_DIR__" || exit 1; while true; do "__BRAINLAYER_BIN__" enrich --mode realtime --limit 200000 --since-hours 87600; sleep 5; done</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__BRAINLAYER_DIR__</string>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/brainlayer-enrichment.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/brainlayer-enrichment.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:__HOME__/.local/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONPATH</key>
        <string>__BRAINLAYER_DIR__/src</string>
        <key>BRAINLAYER_STALL_TIMEOUT</key>
        <string>300</string>
        <key>BRAINLAYER_ENRICH_RATE</key>
        <string>15</string>
        <key>BRAINLAYER_ENRICH_CONCURRENCY</key>
        <string>18</string>
        <key>BRAINLAYER_GEMINI_SERVICE_TIER</key>
        <string>flex</string>
        <key>GOOGLE_API_KEY</key>
        <string>__GOOGLE_API_KEY__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
