# scripts/launchd/com.brainlayer.hotlane-brainbar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.brainlayer.hotlane-brainbar</string>
	<key>ProgramArguments</key>
	<array>
		<string>__BRAINLAYER_ENV_RUN__</string>
		<string>__BRAINLAYER_PYTHON__</string>
		<string>__BRAINLAYER_DIR__/scripts/hotlane_brainbar_daemon.py</string>
		<string>--interval</string>
		<string>1.0</string>
		<string>--recent-limit</string>
		<string>5</string>
		<string>--backlog-interval</string>
		<string>10.0</string>
		<string>--backlog-batch</string>
		<string>4</string>
		<string>--enrich-interval</string>
		<string>10.0</string>
		<string>--enrich-limit</string>
		<string>5</string>
		<string>--enrich-since-hours</string>
		<string>87600</string>
	</array>
	<key>WorkingDirectory</key>
	<string>__BRAINLAYER_DIR__</string>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/brainlayer/hotlane-brainbar.out.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/brainlayer/hotlane-brainbar.err.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.local/bin</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
		<key>BRAINLAYER_ENV_FILE</key>
		<string>__BRAINLAYER_ENV_FILE__</string>
		<key>BRAINLAYER_REQUIRE_GOOGLE_API_KEY</key>
		<string>1</string>
		<key>BRAINLAYER_REPO_ROOT</key>
		<string>__BRAINLAYER_DIR__</string>
		<key>BRAINLAYER_GEMINI_SERVICE_TIER</key>
		<string>flex</string>
		<key>BRAINLAYER_ENRICH_CONCURRENCY</key>
		<string>4</string>
		<key>BRAINLAYER_ENRICH_RATE</key>
		<string>15</string>
		<key>BRAINLAYER_LAUNCHD_SERVICE</key>
		<string>hotlane-brainbar</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>Nice</key>
	<integer>10</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ExitTimeOut</key>
	<integer>30</integer>
	<key>LowPriorityIO</key>
	<true/>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>4096</integer>
	</dict>
</dict>
</plist>
