# scripts/launchd/com.brainlayer.enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brainlayer.enrich</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BRAINLAYER_BIN__</string>
        <string>enrich</string>
        <string>--batch-size</string>
        <string>50</string>
        <string>--max</string>
        <string>500</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/share/brainlayer/logs/enrich.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/share/brainlayer/logs/enrich.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:__HOME__/.local/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>BRAINLAYER_STALL_TIMEOUT</key>
        <string>300</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>15</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
