# scripts/com.zikaron.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zikaron.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/etanheyman/Gits/zikaron/.venv/bin/python</string>
        <string>/Users/etanheyman/Gits/zikaron/scripts/watcher.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/zikaron-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/zikaron-watcher.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/etanheyman</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/etanheyman/Gits/zikaron/.venv/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/etanheyman/Gits/zikaron</string>
</dict>
</plist>
